# deploy/systemd/system/ravenwire-log-prune.service
[Unit]
Description=Prune RavenWire host log handoff files

[Service]
Type=oneshot
ExecStart=/usr/local/libexec/ravenwire-prune-logs
