# quadlet/sensor-pod/pcap-ring-writer.container
[Unit]
Description=PCAP Ring Writer — Alert-Driven Rolling Ring Buffer
PartOf=capture-pipeline.target
After=sensor-agent.service
Requires=sensor-agent.service

[Container]
Image=ghcr.io/sensor-stack/pcap-ring-writer:latest
Network=host
UserNS=keep-id:uid=10000,gid=10000
User=10000:10000
# Drop all capabilities first, then add only CAP_NET_RAW for AF_PACKET socket
DropCapability=all
AddCapability=CAP_NET_RAW
# Memory-mapped ring buffer in shared memory (read-write)
Volume=/dev/shm:/dev/shm:rw
# Unix socket directory shared with sensor_agent PCAP Manager
Volume=/var/run/sensor:/var/run/sensor:rw,z
# BPF filter profile (read-only; Sensor_Agent sends updates via control socket)
Volume=/etc/sensor/bpf_filters.conf:/etc/sensor/bpf_filters.conf:ro,z
# All node-specific values via environment variables — no hardcoded paths
Environment=CAPTURE_IFACE=${CAPTURE_IFACE}
Environment=PCAP_RING_FANOUT_GROUP=${PCAP_RING_FANOUT_GROUP:-4}
Environment=RING_PATH=/dev/shm/sensor_pcap_ring
Environment=RING_SIZE_MB=${RING_SIZE_MB:-4096}
Environment=CONTROL_SOCK=/var/run/sensor/pcap_ring.sock
Exec=/usr/local/bin/pcap-ring-writer
Restart=on-failure

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=capture-pipeline.target
