# deploy/quadlet/management-pod/config-manager.container
[Unit]
Description=Config Manager — Elixir/Phoenix LiveView Management Application
After=network.target

[Container]
Image=localhost/ravenwire/config-manager:test
Network=host
Volume=/data/config_manager:/data:rw
Volume=/data/ca:/data/ca:rw
Volume=/data/metrics:/data/metrics:rw
DropCapability=all
Environment=SENSOR_DB_PATH=/data/config_manager.db
Environment=SENSOR_CA_PATH=/data/ca
Environment=MIX_ENV=dev
Environment=PHX_HOST=0.0.0.0
Environment=SECRET_KEY_BASE=demo_secret_key_base_64chars_long_for_dev_only_not_for_prod_use
Environment=AUTO_ENROLL_FIRST=true
Environment=SENSOR_PORT=4000
Environment=SENSOR_GRPC_PORT=9090

[Service]
Restart=on-failure
RestartSec=10

[Install]
WantedBy=management-pod.target
