# deploy/quadlet/management-pod/management-pod.target
[Unit]
Description=Management Pod — Config Manager and Shared Services
Documentation=https://github.com/ravenwire/ravenwire
Wants=config-manager.service

[Install]
WantedBy=multi-user.target
