# deploy/quadlet/sensor-pod/pcap-ring-writer.container
[Unit]
Description=PCAP Ring Writer — Alert-Driven Rolling Ring Buffer
PartOf=capture-pipeline.target
After=sensor-agent.service
Requires=sensor-agent.service

[Container]
Image=localhost/ravenwire/pcap-ring-writer:test
Network=host
# Drop all capabilities first, then add only CAP_NET_RAW for AF_PACKET socket
DropCapability=all
AddCapability=CAP_NET_RAW
# Memory-mapped ring buffer in shared memory (read-write)
Volume=/dev/shm:/dev/shm:rw
# Unix socket directory shared with sensor_agent PCAP Manager
Volume=/var/run/sensor:/var/run/sensor:rw
# BPF filter profile (read-only; Sensor_Agent sends updates via control socket)
Volume=/etc/sensor/bpf_filters.conf:/etc/sensor/bpf_filters.conf:ro
# All node-specific values via environment variables — no hardcoded paths
Environment=CAPTURE_IFACE=${CAPTURE_IFACE}
Environment=PCAP_RING_FANOUT_GROUP=4
Environment=RING_PATH=/dev/shm/sensor_pcap_ring
Environment=RING_SIZE_MB=512
Environment=CONTROL_SOCK=/var/run/sensor/pcap_ring.sock

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=capture-pipeline.target
