# deploy/quadlet/sensor-pod/suricata.container
[Unit]
Description=Suricata Intrusion Detection System
PartOf=analysis-pipeline.target
After=sensor-agent.service capture-pipeline.target
Requires=sensor-agent.service capture-pipeline.target

[Container]
Image=docker.io/jasonish/suricata:latest
Network=host
PodmanArgs=--entrypoint suricata
# Drop all capabilities first, then add only what AF_PACKET requires
DropCapability=all
AddCapability=CAP_NET_RAW
AddCapability=CAP_NET_ADMIN
# Config: suricata.yaml (read-only; Sensor_Agent writes BPF reload config)
Volume=/etc/sensor/suricata:/etc/suricata:ro
# BPF filter profile (read-only; Sensor_Agent writes updates here)
Volume=/etc/sensor/bpf_filters.conf:/etc/sensor/bpf_filters.conf:ro
# Log handoff path to Vector (read-write)
Volume=/var/sensor/logs/suricata:/logs/suricata:rw
# All node-specific values via environment variables — no hardcoded paths
Environment=CAPTURE_IFACE=${CAPTURE_IFACE}
Environment=SURICATA_FANOUT_GROUP=${SURICATA_FANOUT_GROUP:-2}
Exec=-c /etc/suricata/suricata.yaml --af-packet=${CAPTURE_IFACE} -l /logs/suricata

[Service]
Restart=on-failure
RestartSec=5

[Install]
WantedBy=analysis-pipeline.target
